# upstart/opentribot-boot.conf
start on startup

task
console output

script
 stack_file=/tmp/opentribot-network-interface-stack
 rm -f $stack_file 2>/dev/null || true
end script
