# upstart/opentribot.conf
description	"bringup opentribot"
 
start on filesystem
stop on runlevel [016]
 
console output
respawn

exec opentribot-start
pre-stop exec opentribot-stop
